# deploy/award.conf
[program:award]
command=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 80
directory=/home/user/award-system/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/award.log
stderr_logfile=/var/log/award_err.log
user=root
